# tmp/release/backup/linux_crontab/crontab
0 17 * * * user /home/user/backup/backup.sh >/tmp/backup.log
